# mysql/services/0_mysql/finish
#!/bin/bash
exec mysqladmin -uroot shutdown