# neo4j/services/neo4j/finish
#!/bin/sh
exec kill `ps aux | grep neo4j | grep java | awk '{print $2}'`
