# neo4j/services/neo4j/run
#!/bin/bash

find /opt/neo4j ! -user lo-daemon -exec chown lo-daemon. {} \;
exec chpst -u lo-daemon /opt/neo4j/bin/neo4j console
