# nginx/etc/sv/nginx/run
#!/bin/bash

#chown -R lo-daemon. /var/log/nginx/
#chown -R lo-daemon. /var/lib/nginx/
#chown -R lo-daemon. /var/run/
#chown -R lo-daemon. /etc/nginx/
exec /usr/sbin/nginx 
